# scripts/systemd/ownprem-orchestrator.service
[Unit]
Description=Ownprem Orchestrator
Documentation=https://github.com/ownprem/ownprem
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ownprem
Group=ownprem

# Working directory
WorkingDirectory=/opt/ownprem/repo

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/etc/ownprem/orchestrator.env

# Start command
ExecStart=/usr/bin/node /opt/ownprem/repo/apps/orchestrator/dist/index.js

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/ownprem /var/log/ownprem

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ownprem-orchestrator

[Install]
WantedBy=multi-user.target
