# scripts/systemd/ownprem-agent.service
[Unit]
Description=Ownprem Agent
Documentation=https://github.com/ownprem/ownprem
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ownprem
Group=ownprem

# Working directory
WorkingDirectory=/opt/ownprem/repo

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/etc/ownprem/agent.env

# Start command
ExecStart=/usr/bin/node /opt/ownprem/repo/apps/agent/dist/index.js

# Restart policy
Restart=always
RestartSec=10

# Security hardening (less restrictive - agent needs to manage apps)
NoNewPrivileges=no
PrivateTmp=no
ProtectSystem=no
ProtectHome=yes

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ownprem-agent

[Install]
WantedBy=multi-user.target
